# docker/supervisord-services.ini
[program:nginx]
command=/usr/sbin/nginx -c /home/app/flask/docker/nginx.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes = 0

[program:uwsgi]
command=/usr/sbin/uwsgi --ini /home/app/flask/docker/cla_public.ini --die-on-term
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes = 0
