# docker/supervisord-services.ini
[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
stdout_logfile=/dev/stdout # will only be effective if the program is also logging to stdout
stderr_logfile=/dev/stderr

[program:uwsgi]
command=/usr/sbin/uwsgi --ini /home/app/flask/docker/cla_public.ini --die-on-term
stdout_logfile=/dev/stdout # will only be effective if the program is also logging to stdout
stderr_logfile=/dev/stderr
