# docker/supervisord-services.ini
[program:nginx]
command=/usr/sbin/nginx -c /home/app/flask/docker/nginx.conf

[program:uwsgi]
command=/usr/sbin/uwsgi --ini /home/app/flask/docker/cla_public.ini --die-on-term
